# dev-db/arangodb/files/arangodb.service
# arangodb.service

[Unit]
Description=ArangoDB Database

[Service]
Type=oneshot
#RemainAfterExit=yes
PIDFile=/var/run/arangod.pid
User=root
ExecStart=/usr/sbin/arangod --pid-file /var/run/arangod.pid --server.uid arangodb --daemon
ExecStop=/bin/kill -2 $MAINPID

[Install]
WantedBy=multi-user.target
